# com.openclaw.ops-proxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openclaw.ops-proxy</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>-m</string>
        <string>ops_proxy.cli</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/taasaa/dev/ops/ops-proxy</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>TG_BOT_TOKEN</key>
        <string></string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/taasaa/.openclaw-ops/ops-proxy/ops-proxy.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/taasaa/.openclaw-ops/ops-proxy/ops-proxy.err</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
